# recipes-version-ctrl/hw-githash/files/hw-githash.service
[Unit]
Description=HW Git Hash Service

[Service]
Type=oneshot
ExecStart=/usr/bin/hw-githash

[Install]
WantedBy=multi-user.target
